# greeter.service
[Unit]
Description=My Greeting Service

[Service]
Type=oneshot
ExecStart=/home/joy/Code/linux/greeter.sh