# testconf_watch.path
[Unit]
Description=Path unit for watching for changes in "config"
[Path]
PathModified=/etc/test.conf
Unit=testconf_watch.service
[Install]
WantedBy=multi-user.target
